# deployment/systemd/osmosis-spot.service
[Unit]
Description=Osmosis Spot-Only DEX Node
Documentation=https://github.com/osmosis-labs/osmosis
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=osmosis
Group=osmosis
WorkingDirectory=/home/osmosis
ExecStart=/usr/local/bin/osmosisd start \
    --home /home/osmosis/.osmosisd \
    --log_level info \
    --minimum-gas-prices 0.0025uosmo
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=120
Restart=on-failure
RestartSec=5
StartLimitInterval=0

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/osmosis/.osmosisd
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

# Resource limits
LimitNOFILE=65535
LimitNPROC=32768

# Environment variables for spot-only mode
Environment=SPOT_ONLY_MODE=true
Environment=DISABLE_LEVERAGE_MODULES=true
Environment=GOVERNANCE_SAFEGUARDS_ENABLED=true
Environment=CHAIN_ID=osmosis-spot-1

[Install]
WantedBy=multi-user.target

# Installation instructions:
# 1. Copy this file to /etc/systemd/system/osmosis-spot.service
# 2. Create osmosis user: sudo useradd -r -s /bin/false osmosis
# 3. Create home directory: sudo mkdir -p /home/osmosis/.osmosisd
# 4. Set ownership: sudo chown -R osmosis:osmosis /home/osmosis
# 5. Reload systemd: sudo systemctl daemon-reload
# 6. Enable service: sudo systemctl enable osmosis-spot
# 7. Start service: sudo systemctl start osmosis-spot
# 8. Check status: sudo systemctl status osmosis-spot
# 9. View logs: sudo journalctl -u osmosis-spot -f